# packages/ws-web-aiohttp/ws_web_aiohttp-0.1a52-py3-none-any.whl/ws_web_aiohttp/ws_web_aiohttp.service
[Unit]
Description=web sheets web server

[Service]
ExecStart=/usr/local/bin/ws_web_aiohttp runserver --conf_dir=/etc/ws_web_aiohttp/conf simple

[Install]
WantedBy=multi-user.target

