# packages/greenscreen-control/greenscreen_control-0.0.2.tar.gz/greenscreen_control-0.0.2/greenscreen_control.service
[Unit]
Description=GreenScreen Control server
After=network.target

[Service]
# *** Replace APPID with the AppID for your GreenScreen control App ***
ExecStart=/usr/local/bin/greenscreen_control_server -l INFO -a APPID
Restart=on-abort
User=nobody
Group=nogroup

[Install]
WantedBy=multi-user.target
