# packages/systemdunitextras/systemdunitextras-0.2.6.tar.gz/systemdunitextras-0.2.6/systemdunitextras/data/etc/systemd/system/systemd_fail_handler@.service
[Unit]
Description=Systemd fail handler from systemdunitextras pypi package.
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/systemd_fail_handler -i %I -host %H

[Install]
WantedBy=multi-user.target
